# net-misc/modemmanager/files/init.d.ModemManager
#!/sbin/openrc-run

description="ModemManager Daemon"
pidfile="/run/ModemManager.pid"
command="/usr/sbin/ModemManager"
command_background="yes"

depend() {
	need dbus
	before NetworkManager
}
